# launchagents/com.ibarsi.lms-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ibarsi.lms-server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>export PATH="$HOME/.lmstudio/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; lms server start --port 1234 --bind 127.0.0.1; sleep 5; lms load qwen2.5-coder-7b-instruct-mlx --yes</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
